# docker/s6-rc.d/josh-create-local/up
#!/command/execlineb

define data_dir /data/git

josh-ensure-dir -p ${data_dir}
josh-ensure-owner -p ${data_dir} -o git:git
josh-ensure-mode -p ${data_dir} -m 640
