# docker/s6-rc.d/sshd/run
#!/command/execlineb -P

/usr/sbin/sshd -e -D -h /data/keys/.ssh/id_ed25519
